# voice-assistant.service
[Unit]
Description=ReSpeaker Voice Assistant Service
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
ExecStart=/home/pi/.local/bin/uv run voice-assistant test-realtime
WorkingDirectory=/home/pi/llm-assistant/voice-assistant
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10
User=pi

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/pi/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target

